# hpl-startup.service
[Unit]
Description=Executing HPL on startup
After=network.target
[Service]
ExecStart=/bin/bash -c 'cd /root/hpl-2.3/bin/rpi; /usr/bin/mpirun --allow-run-as-root -np 1 xhpl'
[Install]
WantedBy=default.target
